# systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnStartupSec=5s
OnUnitActiveSec=15m
Unit=mbsync.service

[Install]
WantedBy=timers.target
