# systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=20s
OnUnitActiveSec=15m
Unit=mbsync.service

[Install]
WantedBy=timers.target
